# templates/etc/profile.d/goenv.sh.erb
export GOENV_ROOT="<%= scope.lookupvar('goenv::params::goenv_root') -%>"
[ -d "$GOENV_ROOT" ] && eval "$(goenv init -)"
